# systemd/datang-reader-gui.service
[Unit]
Description=Datang Reader GUI Client
Documentation=https://github.com/yourusername/datang-reader
After=network.target docker.service graphical.target
Wants=docker.service

[Service]
Type=simple
WorkingDirectory=/home/user/Datang-Reader/client
ExecStartPre=/bin/sleep 10
ExecStart=/home/user/Datang-Reader/client/run-gui.sh
Restart=always
RestartSec=10

# X11 Display configuration
Environment="DISPLAY=:0"
Environment="QT_QPA_PLATFORM=xcb"
Environment="XAUTHORITY=/home/user/.Xauthority"

# Python environment
Environment="PYTHONUNBUFFERED=1"

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=datang-reader-gui

[Install]
WantedBy=graphical.target
