# systemd/input-client.service
[Unit]
Description=Datang Reader Input Client (Docker Mode)
Documentation=https://github.com/yourusername/datang-reader
After=network.target docker.service
Requires=docker.service

[Service]
Type=simple
User=root
WorkingDirectory=/opt/datang-reader
ExecStartPre=/bin/sleep 5
ExecStart=/usr/bin/python3 /opt/datang-reader/input_client.py --log-file /var/log/datang-input-client.log

# Restart policy
Restart=always
RestartSec=10

# Environment
Environment="PYTHONUNBUFFERED=1"

# Resource limits (optional)
# MemoryLimit=100M
# CPUQuota=10%

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=datang-input-client

[Install]
WantedBy=multi-user.target
